# device_query/device_query_service.conf
description "labshare device query service"

start on (filesystem)
stop on runlevel [016]

respawn
setuid nobody
setgid nogroup
chdir <device_query_dir>

exec /usr/bin/python3 device_query.py